# etc-conf/dbus/rhsm-facts.service
[Unit]
Description=RHSM system Facts dbus service
After=syslog.target network.target

[Service]
Type=dbus
BusName=com.redhat.RHSM1.Facts
ExecStart=/usr/libexec/rhsm-facts-service

[Install]
WantedBy=basic.target
