# etc-conf/dbus/systemd/rhsm-facts.service.template
[Unit]
Description=RHSM system Facts dbus service
After=syslog.target network.target

[Service]
Type=dbus
BusName=com.redhat.RHSM1.Facts
ExecStart=${libexecdir}/rhsm-facts-service

[Install]
WantedBy=basic.target
